# cron.ini
[program:crond]
command=/usr/bin/crond -n
process_name=%(program_name)s
autostart=true
autorestart=true
priority=10
stopwaitsecs = 3
stopasgroup = false
killasgroup = true
stopsignal=QUIT
user=root
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
